# crontab
*/1 * * * * /home/detect/splunk_detect/detect_wrapper.sh
